# cron/crontab
#crontab

*/30  *  *  *  *  root *default/cron_test
